# orb_listener
### BEGIN INIT INFO
# Provides: Orb listener script
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: BlinkM Powered Orb
# Description: Weather, Mood Lights, Clock
### END INIT INFO


#! /bin/sh
# /etc/init.d/orb_listener


export HOME
case "$1" in
    start)
        echo "Starting Orb"
        /usr/local/bin/orb_listener.py  2>&1 &
    ;;
    stop)
        echo "Stopping Orb"
        LCD_PID=`ps auxwww | grep orb_listener.py | head -1 | awk '{print $2}'`
        kill -9 $LCD_PID
    ;;
    *)
        echo "Usage: /etc/init.d/orb_listener {start|stop}"
        exit 1
    ;;
esac
exit 0
